# Appify.AppDir/usr/share/appify/systemd/appify-updater.service
[Unit]
Description=Appify Updater Service
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/appify-updater.sh
